# script/api.clearspending.ini
[program:api_clearspending_prod_9021]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9021
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stdout.log

[program:api_clearspending_prod_9022]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9022
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stdout.log

[program:api_clearspending_prod_9023]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9023
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stdout.log

[program:api_clearspending_prod_9024]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9024
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stdout.log

[program:api_clearspending_prod_9025]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9025
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stdout.log

[program:api_clearspending_prod_9026]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9026
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stdout.log

[program:api_clearspending_prod_9027]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9027
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stdout.log

[program:api_clearspending_prod_9028]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9028
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stdout.log

[program:api_clearspending_prod_9029]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9029
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-close-stdout.log

[program:api_clearspending_prod_9030]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9030
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-open-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-open-stdout.log

[program:api_clearspending_prod_9031]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9031
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-open-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-open-stdout.log

[program:api_clearspending_prod_9032]
command=/var/www/api.clearspending.ru/cs_api/bin/python /var/www/api.clearspending.ru_v3/api/tornading.py API 9032
directory=/var/www/api.clearspending.ru_v3
environment=PYTHONPATH=/var/www/api.clearspending.ru_v3
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=sdemidenko
stderr_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-open-stderr.log
stdout_logfile = /var/www/api.clearspending.ru_v3/logs/tornado-open-stdout.log

[group:api_clearspending_close]
programs=api_clearspending_prod_9021,api_clearspending_prod_9022,api_clearspending_prod_9023,api_clearspending_prod_9024,api_clearspending_prod_9025,api_clearspending_prod_9026,api_clearspending_prod_9027,api_clearspending_prod_9028,api_clearspending_prod_9029

[group:api_clearspending_open]
programs=api_clearspending_prod_9030,api_clearspending_prod_9031,api_clearspending_prod_9032